# airootfs/etc/profile.d/expand_cowsize.sh
#!/bin/bash
echo 'password' | sudo -S mount -o remount,size=1G /run/archiso/cowspace
